# config/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * cd /var/www/drupal && export $(cat ../container.env | xargs) && drush cron
